# root/etc/systemd/system/php70-php-fpm.service.d/nethserver.conf
[Service]
RuntimeDirectory=php70-php-fpm
RuntimeDirectoryMode=0755
