# root/etc/systemd/system/php56-php-fpm.service.d/nethserver.conf
[Service]
RuntimeDirectory=php56-php-fpm
RuntimeDirectoryMode=0755
