# root/etc/systemd/system/php74-php-fpm.service.d/nethserver.conf
[Service]
RuntimeDirectory=php74-php-fpm
RuntimeDirectoryMode=0755
